# dev/ansible/roles/mw_pg_backup_db/templates/pg_backup.service.j2
[Unit]
Description=Zabbix Database Backup
After=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pg_backup.sh
User=root