# dev/ansible/roles/mw_backup_db/templates/pg_backup.timer.j2
[Unit]
Description=Run pg_backup daily
[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true
[Install]
WantedBy=timers.target