# dev/ansible/roles/web_backup/templates/web_backup.timer.j2
[Unit]
Description=Run web_backup daily
[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true
[Install]
WantedBy=timers.target