# dev/ansible/roles/web_backup/templates/web_backup.timer.j2
[Unit]
Description=Run web_backup каждые 3 дня
[Timer]
OnCalendar=*-*-1/3 03:00:00
Persistent=true
[Install]
WantedBy=timers.target